# deploy/supervisor-backend.conf
; KERNELiOS Django backend — supervisor config for Cloudways
; Copy to /etc/supervisor/conf.d/ and edit paths/env vars

[program:kernelios-backend]
command=/var/www/YOURAPP/new-version/backend/.venv/bin/gunicorn kernelios.wsgi:application --workers 3 --bind 127.0.0.1:8000 --timeout 120 --access-logfile /var/www/YOURAPP/new-version/backend/logs/gunicorn-access.log --error-logfile /var/www/YOURAPP/new-version/backend/logs/gunicorn-error.log
directory=/var/www/YOURAPP/new-version/backend
user=www-data
group=www-data
autostart=true
autorestart=true
redirect_stderr=false
stderr_logfile=/var/log/supervisor/kernelios-backend.err.log
stdout_logfile=/var/log/supervisor/kernelios-backend.out.log
; Set all env vars inline here OR use a .env file read by python-dotenv (already wired)
environment=DJANGO_DEBUG="false",DJANGO_ALLOWED_HOSTS="yourdomain.com"
