# deploy/supervisor-frontend.conf
; KERNELiOS Next.js frontend — supervisor config for Cloudways

[program:kernelios-frontend]
command=/usr/bin/node /var/www/YOURAPP/new-version/frontend/node_modules/.bin/next start --hostname 127.0.0.1 --port 3001
directory=/var/www/YOURAPP/new-version/frontend
user=www-data
group=www-data
autostart=true
autorestart=true
redirect_stderr=false
stderr_logfile=/var/log/supervisor/kernelios-frontend.err.log
stdout_logfile=/var/log/supervisor/kernelios-frontend.out.log
environment=NODE_ENV="production",NEXT_PUBLIC_API_URL="https://yourdomain.com"
